# centos-7/etc/profile.d/00-colors.sh
# Prompt colors constants
readonly NORM="\033[0m"
readonly BOLD="\033[1m"
readonly GREY="\033[38;5;242m"
readonly ACCT="\033[38;5;{{accentcolor}}m"
readonly WHIT="\033[38;5;255m"

# Message boxes size
readonly BOXSIZE=80
